# thirdparty-conf/supervisord/apiaxle-api.conf
[program:apiaxle-api]
process_name = apiaxle-api-%(process_num)s
command = coffee apiaxle_api.coffee %(process_num)s
directory = /home/apiaxle/apiaxle/api
numprocs = 4
numprocs_start = 3000
user=apiaxle
redirect_stderr=true
stdout_logfile=/var/log/apiaxle/api-%(process_num)s-stdout.log
stderr_logfile=/var/log/apiaxle/api-%(process_num)s-stderr.log

