# thirdparty-conf/supervisord/apiaxle-dash.conf
[program:apiaxle-dash]
process_name = apiaxle-dash-%(process_num)s
command = coffee app.coffee %(process_num)s
directory = /home/apiaxle/axle-dash
numprocs = 1
numprocs_start = 8080
user=apiaxle
redirect_stderr=true
stdout_logfile=/var/log/apiaxle/dash-%(process_num)s-stdout.log
stderr_logfile=/var/log/apiaxle/dash-%(process_num)s-stderr.log

