# thirdparty-conf/supervisord/apiaxle-proxy.conf
[program:apiaxle-proxy]
process_name = apiaxle-proxy-%(process_num)s
command = coffee apiaxle_proxy.coffee %(process_num)s
directory = /home/apiaxle/apiaxle/proxy
numprocs = 4
numprocs_start = 2000
user=apiaxle
redirect_stderr=true
stdout_logfile=/var/log/apiaxle/proxy-%(process_num)s-stdout.log
stderr_logfile=/var/log/apiaxle/proxy-%(process_num)s-stderr.log

